# systemd/staticwaves-comfyui.service
[Unit]
Description=StaticWaves ComfyUI Server
After=network.target
Documentation=https://github.com/ssiens-oss/ssiens-oss-static_pod

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ComfyUI
ExecStart=/usr/bin/python3 main.py --listen 0.0.0.0 --port 8188
Restart=always
RestartSec=10

# Environment for GPU
Environment="CUDA_VISIBLE_DEVICES=0"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=staticwaves-comfyui

# Security
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
