# systemd/staticwaves-pod-worker.service
[Unit]
Description=StaticWaves POD ComfyUI Worker
After=network.target
Documentation=https://github.com/ssiens-oss/ssiens-oss-static_pod

[Service]
Type=simple
User=root
WorkingDirectory=/opt/staticwaves-pod
ExecStart=/usr/bin/python3 /opt/staticwaves-pod/workers/comfy_worker.py
Restart=always
RestartSec=10

# Environment
EnvironmentFile=-/opt/staticwaves-pod/config/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=staticwaves-pod-worker

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
